# converted_sigma_rules_to_elastalert/builtin/security/win_security_susp_scheduled_task_update.yml
description: |
  Detects update to a scheduled task event that contain suspicious keywords.
name: Suspicious Scheduled Task Update
index: winlogbeat-*
filter:
- query:
    query_string:
      query: winlog.channel.keyword:security AND (event.code.keyword:4702 AND (winlog.event_data.TaskContentNew.keyword:(*\\appdata\\local\\temp\\*
        OR *\\appdata\\roaming\\* OR *\\users\\public\\* OR *\\windows\\temp\\* OR
        *c\:\\temp\\* OR *\\desktop\\* OR *\\downloads\\* OR *\\temporary\ internet*
        OR *c\:\\programdata\\* OR *c\:\\perflogs\\*)) AND (winlog.event_data.TaskContentNew.keyword:(*regsvr32*
        OR *rundll32* OR *cmd.exe\<\/command\>* OR *cmd\<\/command\>* OR *\<arguments\>\/c\
        * OR *\<arguments\>\/k\ * OR *\<arguments\>\/r\ * OR *powershell* OR *pwsh*
        OR *mshta* OR *wscript* OR *cscript* OR *certutil* OR *bitsadmin* OR *bash.exe*
        OR *bash\ * OR *scrcons* OR *wmic\ * OR *wmic.exe* OR *forfiles* OR *scriptrunner*
        OR *hh.exe*)))
type: any
priority: 3
mitre_attack_tactic:
- execution
- privilege-escalation
- persistence
mitre_attack_technique:
- t1053.005
risk: high
