# converted_sigma_rules_to_elastalert/builtin/taskscheduler/win_taskscheduler_lolbin_execution_via_task_scheduler.yml
description: |
  Detects the execution of Scheduled Tasks where the program being run is located in a suspicious location or where it is an unusual program to be run from a Scheduled Task
name: Scheduled Task Executed Uncommon LOLBIN
index: winlogbeat-*
filter:
- query:
    query_string:
      query: winlog.channel.keyword:microsoft\-windows\-taskscheduler\/operational
        AND (event.code.keyword:129 AND (winlog.event_data.Path.keyword:(*\\calc.exe
        OR *\\cscript.exe OR *\\mshta.exe OR *\\mspaint.exe OR *\\notepad.exe OR *\\regsvr32.exe
        OR *\\wscript.exe)))
type: any
priority: 2
mitre_attack_tactic:
- persistence
mitre_attack_technique:
- t1053.005
risk: medium
