# converted_sigma_rules_to_elastalert/process_creation/proc_creation_win_schtasks_reg_loader_encoded.yml
description: |
  Detects the creation of a schtask that potentially executes a base64 encoded payload stored in the Windows Registry using PowerShell.
name: Scheduled Task Executing Encoded Payload from Registry
index: winlogbeat-*
filter:
- query:
    query_string:
      query: (winlog.event_data.image.keyword:*\\schtasks.exe OR winlog.event_data.originalfilename.keyword:schtasks.exe)
        AND process.command_line.keyword:*\/create* AND (process.command_line.keyword:(*frombase64string*
        OR *encodedcommand*)) AND (process.command_line.keyword:(*get\-itemproperty*
        OR *\ gp\ *)) AND (process.command_line.keyword:(*hkcu\:* OR *hklm\:* OR *registry\:\:*
        OR *hkey_*))
type: any
priority: 3
mitre_attack_tactic:
- execution
- persistence
mitre_attack_technique:
- t1053.005
- t1059.001
risk: high
