# converted_sigma_rules_to_elastalert/process_creation/proc_creation_win_schtasks_creation.yml
description: |
  Detects the creation of scheduled tasks by user accounts via the "schtasks" utility.
name: Scheduled Task Creation Via Schtasks.EXE
index: winlogbeat-*
filter:
- query:
    query_string:
      query: (process.executable.keyword:*\\schtasks.exe AND process.command_line.keyword:*\
        \/create\ *) AND (NOT (user.name.keyword:(*authori* OR *autori*)))
type: any
priority: 1
mitre_attack_tactic:
- execution
- persistence
- privilege-escalation
mitre_attack_technique:
- t1053.005
- s0111
- car.2013-08-001
- stp.1u
risk: low
